# deploy/gusto-admin.service
[Unit]
Description=Gusto top-10 admin application
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=gusto-cantgohome
Group=gusto-cantgohome
EnvironmentFile=/etc/gusto-cantgohome/app.env
WorkingDirectory=/opt/gusto-cantgohome/current
ExecStart=/opt/gusto-cantgohome/venv/bin/uvicorn app.admin_app:app --host 127.0.0.1 --port 8011
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/gusto-cantgohome

[Install]
WantedBy=multi-user.target
